# sbin/route_xdcdb-fos.crontab
*/30 0-3,12-23 * * * /bin/bash -l -c "/soft/warehouse-apps-1.0/Manage-XDCDB-Users/PROD/sbin/route_xdcdb-fos.rc start"
